# fok3/scripts/crontab
0 4 * * * /hda5/mud/vme1.1/scripts/backup

